# docker-cron/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/go/bin
* * * * * root go run /hello.go>>/var/log/cron.log 2>&1

